# supervisor/conf.d/spamd.conf
[program:spamd]
process_name = spamd
command = /usr/sbin/spamd --max-children 5
stdout_logfile = /var/log/%(program_name)s.log
stderr_logfile = /var/log/%(program_name)s.log
autorestart = true
